# deploy/nanokvm-cleanup.timer
[Unit]
Description=Run NanoKVM ISO cleanup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
